# platform/broadcom/sonic-platform-modules-h3c/s6850-48y8c-w1/service/sonic-platform-h3c.service
#sonic-platform-h3c
#

[Unit]
Description=sonic-platfrom-h3c kernel restore service
After=networking.service
Before=pmon.service

[Service]
Type=oneshot
ExecStart=-/bin/bash /usr/local/bin/h3ckoinstall

[Install]
WantedBy=multi-user.target

